# archive/cleanup-2026-06-11/temp-clone/scripts/siyang-webhook.service
[Unit]
Description=Siyang CMS Webhook Listener
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/siyang/public
ExecStart=/usr/bin/node /var/www/siyang/public/scripts/webhook-listener.js
Restart=on-failure
RestartSec=10
Environment=WEBHOOK_PORT=3099
Environment=WEBHOOK_SECRET=change-me-in-production
Environment=GIT_REPO_PATH=/var/www/siyang/public
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
